# priceweb-worker.service
[Unit]
Description=PriceWeb Data Worker
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/priceweb_new
# IMPORTANT: Put your real Telegram tokens here if not using env vars!
Environment="TG_BOT_TOKEN=ВАШ_ТОКЕН"
Environment="TG_CHAT_ID=ВАШ_ID"
Environment="PRICE_DB_PATH=/opt/priceweb_new/data/priceweb.db"
Environment="TZ=Europe/Moscow"
ExecStart=/opt/priceweb_new/venv/bin/python /opt/priceweb_new/worker.py
StandardOutput=append:/opt/priceweb_new/cron_log.log
StandardError=append:/opt/priceweb_new/cron_log.log

[Install]
WantedBy=multi-user.target
